# common/files/etc/profile.d/debcomp-userconf.sh
[ -e ~/.debcomp/userconf/installed/common ] || debcomp userconf -f common
